# scripts/tar1090-svelte-history.service
[Unit]
Description=tar1090-svelte History Daemon
Documentation=https://github.com/authrequest/tar1090-svelte
After=network.target

[Service]
Type=simple
User=tar1090
Group=tar1090
WorkingDirectory=/var/lib/tar1090-svelte/history
ExecStart=/opt/tar1090-svelte/scripts/history-daemon.sh
Restart=always
RestartSec=5
Environment="SERVICE_NAME=tar1090-svelte-history"
Environment="RUN_DIR=/var/lib/tar1090-svelte/history"
Environment="SRC_URL=http://127.0.0.1:8080/data/aircraft.json"
Environment="INTERVAL=8"
Environment="HISTORY_SIZE=450"
Environment="CHUNK_SIZE=60"
Environment="PTRACKS=8"
Environment="GZIP_LVL=1"
Environment="ENABLE_COMPRESSION=yes"

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tar1090-svelte/history
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
